# Linux/Code/SystemD/fork_and_keep_parent_process_abnormal/myservice.service
[Unit]
Description=My SystemD Service

[Service]
Type=forking
ExecStart=/usr/bin/myservice
TimeoutStartSec=5
WatchdogSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target